# config/systemd/user/hyprpolkitagent.service.d/override.conf
[Unit]
After=
After=dbus.service graphical-session.target
PartOf=graphical-session.target

[Install]
WantedBy=graphical-session.target